# templates/dockerfiles/ambari/centos7/root_agent/etc/supervisor/conf.d/agent.conf
[program:ambari-agent]
command=/usr/bin/start-ambari-agent
redirect_stderr=true
