# templates/dockerfiles/ambari/centos7/root_server/etc/supervisor/conf.d/server.conf
[program:ambari-server]
command=/usr/bin/start-ambari-server
autorestart=true
stopwaitsecs=60
